# centos/mongo
#!/bin/sh
#chkconfig: 2345 80 80 
#description: auto start mongo

mongod -dbpath "/data/mongo/data"