# centos/sudoku-game
#!/bin/sh
#chkconfig: 2345 80 80 
#description: auto start sudoku game

node /data/web/sudoku-game/bin/startup.js