# centos/sudoku-game
#!/bin/sh
#chkconfig: 2345 80 80 
#description: auto start sudoku game

cd /data/web/sudoku-game
node ./bin/startup.js